# services/nestcam.service
[Unit]
Description=NestCam DIY camera daemon
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
DynamicUser=yes
SupplementaryGroups=video
WorkingDirectory=/opt/nestcam
ExecStartPre=/bin/sleep 5
ExecStart=/usr/bin/python3 /opt/nestcam/nestcam_daemon.py
Restart=always
RestartSec=3
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/etc/nestcam/nestcam.env
StateDirectory=nestcam
RuntimeDirectory=nestcam
TimeoutStopSec=5
KillMode=mixed
SendSIGKILL=yes
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=full
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
