# services/nestcam.service
[Unit]
Description=NestCam DIY camera daemon
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
DynamicUser=yes
SupplementaryGroups=video
WorkingDirectory=/opt/nestcam
ExecStart=/usr/bin/python3 /opt/nestcam/nestcam_daemon.py
Restart=on-failure
RestartSec=3
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/etc/nestcam/nestcam.env
StateDirectory=nestcam
RuntimeDirectory=nestcam
TimeoutStopSec=20
KillMode=mixed
SendSIGKILL=yes
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=full
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
